# docker/service/update/finish
#!/bin/sh

echo "[update] exited."

# if it stops, make sure it won't get automatically started back up again
# by the runsvdir.
sv stop update
