# docker/service/update/finish
#!/bin/sh

SVC="${PWD##*/}"

echo "[$SVC] exited."

set -x

# if it stops, make sure it won't get automatically started back up again
# by the runsvdir.
sv stop $SVC
