# docker/service/app/finish
#!/bin/sh

[[ ! -z "$VERBOSE" ]] && set -x

SVC="${PWD##*/}"
echo "[$SVC] stopped."

# if it stops, make sure it won't get automatically started back up again
# by the runsvdir.
sv stop $SVC
